# .claude/worktrees/agent-a93ab319/systemd/media-control-daemon.service
[Unit]
Description=Media Control Daemon for Hyprland
Documentation=https://github.com/rain/media-control
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=HYPRLAND_INSTANCE_SIGNATURE

[Service]
Type=simple
ExecStart=%h/.cargo/bin/media-control-daemon foreground
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=media_control=info

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t

[Install]
WantedBy=hyprland-session.target
