# .claude/worktrees/agent-a272bd19/systemd/media-control-daemon.socket
[Unit]
Description=Media Control Daemon Socket
PartOf=media-control-daemon.service

[Socket]
ListenStream=%t/media-control-daemon.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
